# telegram-bot.service
[Unit]
Description=Assistente Acadêmico Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=uach
WorkingDirectory=/home/uach/atendimento_alunos_bot
ExecStart=/home/uach/atendimento_alunos_bot/start_rp4.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
